# docker/node/supervisord/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/supervisord.log
pidfile=/var/run/supervisord.pid

[program:node]
directory=/var/www/echo
command=node app.js
autostart=true
autorestart=true
